# etc/systemd/system/tts-server.service
[Unit]
Description=Local TTS server for Home Assistant
After=network.target

[Service]
Type=simple
User=herrfly
WorkingDirectory=/home/herrfly/Desktop/our-home-assistant-main
# Falls du flask per --user installiert hast, setze PATH so, dass ~/.local/bin gefunden wird:
Environment=PATH=/home/herrfly/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/python3 /home/herrfly/Desktop/our-home-assistant-main/tts_server.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
